# minimdnsd.service
[Unit]
Description=minimdnsd server

[Service]
Type=simple
ExecStart=minimdnsd -r


[Install]
WantedBy=multi-user.target
